# scripts/launchd/com.roddy.ccm-aggregator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.roddy.ccm-aggregator</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>__HOME__/claude-code-monitor-multi/dist/bin/ccm.js</string>
    <string>aggregator</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>__HOME__/.claude-monitor/logs/aggregator.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/.claude-monitor/logs/aggregator.err.log</string>

  <key>WorkingDirectory</key>
  <string>__HOME__/claude-code-monitor-multi</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
